# docker/openlava/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command = /usr/sbin/sshd -D
logfile = /var/log/supervisor/sshd.log
autorestart=true

[program:openlava]
command=service openlava restart
